# sample-scripts/bin/alfresco-services.service
[Unit]
Description=Alfresco Content Repository
After=docker.service
Wants=network-online.target docker.socket
Requires=docker.socket

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/alfresco/workspace/docker-alfresco-community
Environment="PATH=/root/.local/bin:/root/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/home/alfresco/workspace/bin/start-stop-service.sh start

[Install]
WantedBy=multi-user.target
